# infra/supervisord.conf
[unix_http_server]
file=/tmp/gpu-inference-supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/gpu-inference/supervisord.log
pidfile=/tmp/gpu-inference-supervisord.pid
childlogdir=/var/log/gpu-inference

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/gpu-inference-supervisor.sock

[program:redis]
command=/bin/bash -lc 'set -a; . /app/.env; set +a; exec redis-server --bind 127.0.0.1 --port "${REDIS_PORT:-6379}" --save "" --appendonly no'
priority=10
autostart=true
autorestart=true
startsecs=2
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/gpu-inference/redis.log
stderr_logfile=/var/log/gpu-inference/redis.err.log

[program:vllm]
command=uv run --env-file /app/.env python -m backend.scripts.run_vllm
directory=/app
priority=20
autostart=true
autorestart=true
startsecs=10
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=180
stdout_logfile=/var/log/gpu-inference/vllm.log
stderr_logfile=/var/log/gpu-inference/vllm.err.log

[program:api]
command=/bin/bash -lc 'set -a; . /app/.env; set +a; uv run --env-file /app/.env python -m backend.scripts.wait_for_vllm && exec uv run --env-file /app/.env uvicorn backend.main:app --host "${APP_HOST:-127.0.0.1}" --port "${APP_PORT:-8000}"'
directory=/app
priority=30
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=60
stdout_logfile=/var/log/gpu-inference/api.log
stderr_logfile=/var/log/gpu-inference/api.err.log

[program:batch-worker]
command=/bin/bash -lc 'uv run --env-file /app/.env python -m backend.scripts.wait_for_vllm && exec uv run --env-file /app/.env python -m backend.workers.batch_worker'
directory=/app
priority=40
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=60
stdout_logfile=/var/log/gpu-inference/batch-worker.log
stderr_logfile=/var/log/gpu-inference/batch-worker.err.log

[program:recovery-scanner]
command=uv run --env-file /app/.env python -m backend.workers.recovery_scanner
directory=/app
priority=41
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=60
stdout_logfile=/var/log/gpu-inference/recovery-scanner.log
stderr_logfile=/var/log/gpu-inference/recovery-scanner.err.log

[program:analytics-flusher]
command=uv run --env-file /app/.env python -m backend.workers.analytics_flusher
directory=/app
priority=42
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=60
stdout_logfile=/var/log/gpu-inference/analytics-flusher.log
stderr_logfile=/var/log/gpu-inference/analytics-flusher.err.log

[program:cloudflared]
command=/bin/bash -lc 'set -a; . /app/.env; set +a; uv run --env-file /app/.env python -m backend.scripts.wait_for_http "http://${APP_HOST:-127.0.0.1}:${APP_PORT:-8000}/health" --timeout 120 && if [ -s /etc/cloudflared/token ]; then exec cloudflared tunnel --no-autoupdate run --token-file /etc/cloudflared/token; else exec cloudflared tunnel --config /etc/cloudflared/config.yml run gpu-inference-backend; fi'
directory=/app
priority=50
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=60
stdout_logfile=/var/log/gpu-inference/cloudflared.log
stderr_logfile=/var/log/gpu-inference/cloudflared.err.log

[program:dcgm-exporter]
command=/bin/bash -lc 'command -v dcgm-exporter >/dev/null 2>&1 && exec dcgm-exporter || { echo "dcgm-exporter not installed; install NVIDIA DCGM exporter before enabling GPU metrics"; exit 0; }'
priority=60
autostart=true
autorestart=false
startsecs=0
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=60
stdout_logfile=/var/log/gpu-inference/dcgm-exporter.log
stderr_logfile=/var/log/gpu-inference/dcgm-exporter.err.log
